# laptop/mic-key.service
[Unit]
Description=Mic LED sync daemon (ESP32 LED follows PA mute state)
After=graphical-session.target

[Service]
ExecStart=/home/silent/Documents/Computarias/Fin-ESP/laptop/mic_key_daemon.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
